# com.archangel.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.archangel.dashboard</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/jordan/Claude/Archangel/marketing-dashboard/venv/bin/python</string>
        <string>/Users/jordan/Claude/Archangel/marketing-dashboard/app.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jordan/Claude/Archangel/marketing-dashboard</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/jordan/Claude/Archangel/marketing-dashboard/dashboard.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jordan/Claude/Archangel/marketing-dashboard/dashboard.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>FLASK_ENV</key>
        <string>production</string>
    </dict>
</dict>
</plist>
